# files/autostart/tilix-terminal.desktop
[Desktop Entry]
Type=Application
Exec=tilix --full-screen
Hidden=false
NoDisplay=false
Name=Terminal
Comment=Start Terminal On Startup
